# meta-rootfsenc/recipes-core/mount-bootimage/files/boot.mount
[Unit]
Before=local-fs.target

[Mount]
What=/dev/disk/by-partlabel/APP
Where=/boot
Type=ext4
